# recipes-connectivity/samba/samba/winbind.service
[Unit]
Description=Samba Winbind Daemon
After=syslog.target network.target nmb.service

[Service]
PIDFile=/var/run/winbindd.pid
ExecStart=@SBINDIR@/winbindd -F
ExecReload=@BASE_BINDIR@/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
